# mylora.service
[Unit]
Description=LoRA database web interface
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/MyLora
ExecStart=/opt/MyLora/venv/bin/python /opt/MyLora/main.py
Restart=always

[Install]
WantedBy=multi-user.target

